# vibelearn-daily-sync.timer
[Unit]
Description=VibeLearn Daily Sync Timer
Requires=vibelearn-daily-sync.service

[Timer]
OnCalendar=*-*-* 00:01:00
Persistent=true

[Install]
WantedBy=timers.target
